# setup-macos/com.local.KeyRemapping.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
    hidutil kernel-level key remapping. Runs once at login, no background process.
    Serves as BT fallback — Karabiner can't intercept Bluetooth HID devices.
    Karabiner overrides these for dongle/USB devices it can see.

    Keychron V6 Max (ANSI) + Lithuanian layout mappings:
      caps_lock (0x39) -> F18 (0x6D)        — used by Hammerspoon as hyper key
      grave_accent (0x35) -> non_us_backslash (0x64) — fixes tilde on ANSI + Lithuanian
-->
<dict>
    <key>Label</key>
    <string>com.local.KeyRemapping</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/hidutil</string>
        <string>property</string>
        <string>--set</string>
        <string>{"UserKeyMapping":[{"HIDKeyboardModifierMappingSrc":0x700000039,"HIDKeyboardModifierMappingDst":0x70000006D},{"HIDKeyboardModifierMappingSrc":0x700000035,"HIDKeyboardModifierMappingDst":0x700000064}]}</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
